# service/etc/systemd/system/manifest.service
[Unit]
Description = prepare works for BUDA IIIF reader

[Service]
ExecStart = /usr/local/bin/manifest.sh start

[Install]
WantedBy = default.target